# ha-addon_seatable/rootfs/etc/services.d/seatable/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Seatable service
# ==============================================================================

bashio::log.info "Starting Seatable..."
exec /usr/bin/seatable.sh